# supervisord.conf
; Example config file for running logserver under supervisord.

[program:logserver]
command=logserver 3042 %(ENV_HOME)s/log

[program:logmd]
command=logmd %(ENV_HOME)s/log
